# deploy/home-podcast.service
[Unit]
Description=Home Podcast Metadata API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=home-podcast
Group=home-podcast
EnvironmentFile=/etc/home-podcast.env
WorkingDirectory=/opt/home-podcast
ExecStart=/opt/home-podcast/home-podcast
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/srv/home-podcast/audio
RuntimeDirectory=home-podcast
RuntimeDirectoryMode=0750
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
